# supv/supv_danmu.conf
[program:danmuG]
command=python3 danmu.py 
user=root
directory=~/quin_danmu/
autostart=true
autorestart=true
startsecs=15
stderr_logfile=~/quin_danmu/log_err.txt
