# supv/supv_danmu.conf
[program:danmuG]
command=python3 danmu.py 
user=root
directory=~/danmuG/
autostart=true
autorestart=true
startsecs=10
stdout_logfile=~/danmuG/log_out.txt
stderr_logfile=~/danmuG/log_err.txt
